# .config/systemd/user/autostart-script.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/home/dle/.config/hypr/scripts/autostart.sh

