# .config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland startup
StopWhenUnneeded=yes
Wants=xdg-desktop-autostart.target
