# .config/systemd/user/suspend.target
[Unit]
Description=User level sleep target
StopWhenUnneeded=yes
